# http.service
[Unit]
Description=Legacy Python HTTP server (deprecated - use Apache/nginx instead)
After=multi-user.target

[Service]
Type=simple
Restart=always
WorkingDirectory=/home/ubuntu/fermmon
ExecStart=/usr/bin/python -m http.server

[Install]
WantedBy=multi-user.target
